# paper_search-main/preprocess.conf
[program: pre_processing]
directory=/root/paper_search
command=/root/miniconda3/bin/python pre_processing.py
autostart=true
stdout_logfile=/root/pre_processing_out.log
